# climateiq/rootfs/etc/services.d/redis/run
#!/usr/bin/execlineb -P
# S6 service: Redis
# Embedded Redis for single-container HA add-on deployment

fdmove -c 2 1

s6-echo "Starting Redis..."

exec redis-server /etc/redis-climateiq.conf
