# ha-addon/rootfs/etc/services.d/nginx/run
#!/usr/bin/execlineb -P
# S6 service: Nginx ingress proxy
# Proxies HA ingress requests to the ClimateIQ backend

fdmove -c 2 1

s6-echo "Starting Nginx ingress proxy..."

# Wait for backend to be available before starting nginx
foreground {
    backtick -n READY {
        loopwhilex -x 0
        backtick -n CODE { curl -sf -o /dev/null -w "%{http_code}" http://127.0.0.1:8420/health }
        importas CODE CODE
        if { s6-test ${CODE} = "200" }
        exit 0
    }
}

exec nginx -g "daemon off;"
