# climateiq/rootfs/etc/services.d/postgresql/run
#!/usr/bin/execlineb -P
# S6 service: PostgreSQL
# Embedded PostgreSQL for single-container HA add-on deployment

fdmove -c 2 1

s6-echo "Starting PostgreSQL..."

# Ensure data directory exists and has correct ownership
if { s6-test -d /data/postgresql/data }

exec su-exec postgres postgres -D /data/postgresql/data
